# atlas-edge-attendance.service
[Unit]
Description=Atlas Edge Attendance Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/limitify/atlas-edge
ExecStart=/usr/bin/python3 /home/limitify/atlas-edge/attendance_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
